# aircon-control/rootfs/etc/services.d/apache/finish
#!/usr/bin/with-contenv bashio
# Finish script for apache service - no cleanup required, exit cleanly
exit 0
